# backend/deploy/templates/todoapi.service.j2
[Unit]
Description=To-Do list backend
After=network.target

[Service]
EnvironmentFile=/etc/todoapp/.env
ExecStart={{ bin_location }}
Restart=always

[Install]
WantedBy=multi-user.target